# build/better-go-playground.service
[Unit]
Description=Better Go Playground
After=network.target
StartLimitIntervalSec=10
StartLimitBurst=5

[Service]
Type=simple
User=go-playground
WorkingDirectory=/opt/better-go-playground
ExecStart=/opt/better-go-playground/playground -f /opt/better-go-playground/data/packages.json
Restart=on-failure
LimitNOFILE=1048576
RestartSec=5

[Install]
WantedBy=multi-user.target
