# build/go-playground.service
[Unit]
Description=go-playground
After=network.target

[Service]
Type=notify
User=go-playground
WorkingDirectory=/home/go-playground
Environment=GOROOT=/usr/local/go
ExecStart=/home/go-playground/target/playground -f ./data/packages.json -debug
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
LimitNOFILE=1048576
NotifyAccess=main

[Install]
WantedBy=multi-user.target
